# scripts/sync-agent-config.service
[Unit]
Description=Sync Codex agent config snapshot

[Service]
Type=oneshot
ExecStart=%h/.codex/bin/sync-agent-config.sh
